# gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn SteamFreeGame:app -b localhost:8080
directory=/deploy/app
