# gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn SteamFreeGame:app -b localhost:5555
directory=/deploy/app